# recipes-support/x11-common/files/xserver-nodm.service.in
[Unit]
Description=Xserver startup without a display manager
StartLimitIntervalSec=60s
StartLimitBurst=10

[Service]
EnvironmentFile=/etc/default/xserver-nodm
User=@USER@
ExecStart=/etc/xserver-nodm/Xserver
Restart=on-failure
RestartSec=3s

[Install]
Alias=display-manager.service
